# supervisor/cems-worker.conf
; supervisor/cems-worker.conf
; CEMS-MY Queue Worker Configuration
; Copy to /etc/supervisor/conf.d/cems-worker.conf on production
; 
; Installation:
;   sudo cp supervisor/cems-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start cems-worker:*
;
; Monitoring:
;   sudo supervisorctl status cems-worker:*
;   sudo tail -f /var/log/supervisor/cems-worker.log

[program:cems-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/cems-my/artisan queue:work redis --sleep=3 --tries=3 --timeout=3600 --max-jobs=1000 --max-time=3600 --queue=high,default,low
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cems-worker.log
stopwaitsecs=3600
startretries=3
priority=100

; Environment variables
environment=HOME="/var/www",USER="www-data"

; Resource limits
; Ensure PHP memory limit is set appropriately in php.ini or via -d flag
