# supervisor/cems-horizon.conf
; supervisor/cems-horizon.conf
; CEMS-MY Laravel Horizon Configuration
; Copy to /etc/supervisor/conf.d/cems-horizon.conf on production
;
; Installation:
;   sudo cp supervisor/cems-horizon.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start cems-horizon
;
; Monitoring:
;   sudo supervisorctl status cems-horizon
;   sudo tail -f /var/log/supervisor/cems-horizon.log

[program:cems-horizon]
process_name=%(program_name)s
command=/usr/bin/php /var/www/cems-my/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cems-horizon.log
stopwaitsecs=3600
startretries=3
priority=200

; Environment variables
environment=HOME="/var/www",USER="www-data"

; Notes:
; - Horizon manages its own worker processes dynamically
; - The 'stopwaitsecs' must be >= longest job timeout (3600 seconds)
; - Horizon dashboard available at /horizon (requires admin role)
; - Ensure APP_ENV=production and QUEUE_CONNECTION=redis in .env
